# test/systemd/bifrost.service
[Unit]
Description=Bifrost Gateway Service (Test)
After=network.target

[Service]
# Use forking type - Master daemonizes and writes PID file before parent exits
Type=forking

# Master writes PID file after successful initialization
PIDFile=/app/logs/bifrost.pid

Environment=GOTRACEBACK=crash
WorkingDirectory=/app

# Validate config before starting
ExecStartPre=/app/bifrost -t -c=conf/config.yaml

# Start in Master-Worker + Daemon mode
# -d: Daemon mode (fork to background)
ExecStart=/app/bifrost -d -c=conf/config.yaml

# Hot reload: send SIGHUP to Master (triggers Worker replacement)
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown: send SIGTERM to Master
ExecStop=/bin/kill -TERM $MAINPID

# Restart on failure
Restart=on-failure

# Allow 30 seconds for graceful shutdown
TimeoutStopSec=30s
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
